# weintegrity.service
[Unit]
Description=WEintegrity Backend API Server
Documentation=https://github.com/your-repo
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/weintegrity
ExecStart=/usr/bin/node /var/www/weintegrity/backend/server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=weintegrity-backend

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/var/www/weintegrity/backend/.env

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
